# etc/supervisor.mongod.conf
[program:mongod]
command=/usr/local/bin/mongod --dbpath /opt/mongodb-data/
autostart=true
startsecs=3
autorestart=true